# ops/systemd/oracle-backend.service
[Unit]
Description=Oracle Site Flask API
After=network.target

[Service]
WorkingDirectory=/opt/oracle-site/backend
EnvironmentFile=/opt/oracle-site/backend/.env
ExecStart=/opt/oracle-site/backend/.venv/bin/gunicorn "app.main:create_app()" --bind 127.0.0.1:8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

